# cloudflare/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for PflegeNavigator
Documentation=https://developers.cloudflare.com/cloudflare-one/connections/connect-apps
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/cloudflared tunnel run --config /data/.openclaw/workspace/cloudflare/config.yml
Restart=always
RestartSec=5
KillMode=process
StandardOutput=append:/var/log/cloudflared/tunnel.log
StandardError=append:/var/log/cloudflared/tunnel.log

[Install]
WantedBy=multi-user.target
